# lib/generators/generators/capistrano/rails/services/templates/puma.service.erb
[Unit]
Description=<%= description %>
After=network.target

[Service]
# Foreground process (do not use --daemon in ExecStart or config.rb)
Type=simple

User=<%= role.user %>
Group=<%= role.user %>

# Specify the path to your puma application root
WorkingDirectory=<%= current_path %>

# Helpful for debugging socket activation, etc.
# Environment=PUMA_DEBUG=1

ExecStart=<%= rvm_do "bundle exec puma -C #{fetch(:puma_conf)}", role %> 
ExecStop=<%= rvm_do "bundle exec pumactl -S #{fetch(:puma_state)} stop", role %>


PIDFile=<%= "#{shared_path}/tmp/pids/puma.pid" %>

Restart=always


[Install]
WantedBy=multi-user.target